# images/wn/supervisord.d/condor.conf
[program:condor]
command=condor_master -f 
autorestart=true
